# supervisor/openbox.ini
[program:openbox]
command=/usr/bin/openbox-session
environment=DISPLAY=":99"
stdout_logfile=/var/log/shared/supervisor.openbox.log
priority=10
