# supervisor/selenium.ini
[program:selenium]
command=java -jar /tmp/selenium-server-standalone.jar
environment=DISPLAY=":99",http_proxy="",https_proxy="",HTTP_PROXY="",HTTPS_PROXY=""
stdout_logfile=/var/log/shared/supervisor.selenium.log
priority=20