# supervisor/twm.ini
[program:twm]
command=twm -display :99
environment=DISPLAY=":99"
stdout_logfile=/var/log/shared/supervisor.twm.log
priority=10
