# supervisor/x11vnc.ini
[program:x11vnc]
command=x11vnc -rfbauth /tmp/vncpass -display :99 -forever
stdout_logfile=/var/log/shared/supervisor.x11vnc.log
priority=20
