# supervisor/xvfb.ini
[program:xvfb]
command=/usr/bin/Xvfb +extension RANDR :99 -ac -screen 0 1280x1024x24
stdout_logfile=/var/log/shared/supervisor.xvfb.log
priority=1
